# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:openclaw]
command=openclaw gateway run --force --port %(ENV_OPENCLAW_GATEWAY_PORT)s
directory=/data/openclaw
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/data/logs/openclaw.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/openclaw-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="/data/openclaw"

[program:opencode]
command=opencode serve --port %(ENV_OPENCODE_PORT)s --hostname 127.0.0.1 --pure
directory=/data/code
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/data/logs/opencode.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/opencode-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="/data/opencode"

[program:langgraph]
command=/opt/langgraph/.venv/bin/langgraph dev --no-reload --port %(ENV_LANGGRAPH_PORT)s --host 0.0.0.0
directory=/data/langgraph
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/data/logs/langgraph.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/data/logs/langgraph-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=HOME="/data/langgraph",PYTHONPATH="/data/langgraph/src"
